# usr/lib/opkg/info/mibrain-service.list
/etc/init.d/mibrain_service
/usr/bin/mibrain_service
