# usr/lib/opkg/info/mico-upgrade.list
/etc/init.d/silentboot
/bin/boardupgrade.sh
/bin/silentboot.sh
/bin/flash.sh
