# usr/lib/opkg/info/mico-ai-crontab.list
/usr/bin/mico_ai_crontab
/etc/init.d/mico_ai_crontab
